# tvheadend/rootfs/etc/s6-overlay/s6-rc.d/tvh-proxy/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TVHeadend
# Starts TVHeadend
# ==============================================================================

bashio::log.info "Starting tvhProxy..."

# Configure app
export TVH_BINDADDR=$(bashio::config 'TVH_BINDADDR')
export TVH_PROXY_HOST=$(bashio::config 'TVH_PROXY_HOST')
export TVH_PROXY_PORT=$(bashio::config 'TVH_PROXY_PORT')
export TVH_TUNER_COUNT=$(bashio::config 'TVH_TUNER_COUNT')
export TVH_WEIGHT=$(bashio::config 'TVH_WEIGHT')
export TVH_CHUNK_SIZE=$(bashio::config 'TVH_CHUNK_SIZE')
export TVH_PROFILE=$(bashio::config 'TVH_PROFILE')

{
    printf "%s\n" "TVH_BINDADDR=\"${TVH_BINDADDR}\""
    printf "%s\n" "TVH_PROXY_HOST=\"${TVH_PROXY_HOST}\""
    printf "%s\n" "TVH_PROXY_PORT=\"${TVH_PROXY_PORT}\""
    printf "%s\n" "TVH_TUNER_COUNT=\"${TVH_TUNER_COUNT}\""
    printf "%s\n" "TVH_WEIGHT=\"${TVH_WEIGHT}\""
    printf "%s\n" "TVH_CHUNK_SIZE=\"${TVH_CHUNK_SIZE}\""
    printf "%s\n" "TVH_PROFILE=\"${TVH_PROFILE}\""
} >> ~/.bashrc

exec python /usr/srv/app/tvhProxy.py
