# tvheadend/rootfs/etc/cron.d/webgrabplus
0 0 * * * /homeassistant/tvheadend/wg++/run.sh
