# tvheadend/rootfs/etc/s6-overlay/s6-rc.d/oscam/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: TVHeadend
# Starts TVHeadend
# ==============================================================================

bashio::log.info "Starting OSCam..."

# Ensure config directory exists
if ! bashio::fs.directory_exists '/homeassistant/oscam/'; then
    bashio::log.info "Creating default configuration directory at /homeassistant/oscam/"
    mkdir -p /homeassistant/oscam
    cp /etc/oscam/* /homeassistant/oscam/
fi

exec /usr/bin/oscam -c /homeassistant/oscam
